# extras/systemd/uconsole-speaker-amp.service
[Unit]
Description=Enable uConsole Speaker Amplifier via sysfs (GPIO 523)
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c '\
  GPIO=523; BASE=/sys/class/gpio; \
  [ -e "$BASE/gpio$GPIO" ] || echo $GPIO > "$BASE/export"; \
  echo out > "$BASE/gpio$GPIO/direction"; \
  echo 1   > "$BASE/gpio$GPIO/value" \
'

[Install]
WantedBy=multi-user.target
